# release/routing-installer/payload/manager-overlay/systemd/xray-router.service
[Unit]
Description=Local Xray router-core for commercial routing
After=network.target nss-lookup.target

[Service]
User=root
ExecStart=/usr/bin/xray run -config /etc/xray-router/config.json
Restart=always
RestartSec=3
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
