# release/antishare-installer/payload/manager-overlay/systemd/hiddify-anti-share.timer
[Unit]
Description=Hiddify Anti-Share cycle timer (every 2 minutes)
Requires=hiddify-anti-share.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
